# files/etc/uci-defaults/99-disable-hnat-detect-rxppd
#!/bin/sh
# Disable broken hnat-detect rxppd auto bridge logic for Airpi wwan0_1 scenario

rm -f /etc/hotplug.d/iface/99-hnat-detect 2>/dev/null || true
rm -rf /usr/share/ucode/hnat 2>/dev/null || true

# 清理运行态 rxppd，避免它继续挂在 br-lan 上
echo 0 > /sys/kernel/debug/hnat/hook_toggle 2>/dev/null || true
ip link set rxppd nomaster 2>/dev/null || true
ip link delete rxppd 2>/dev/null || true

exit 0
